# srcs/mysql/srcs/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:mysql]
command=/usr/bin/mysqld_safe --datadir="/var/lib/mysql/"
autostart=true
autorestart=false # To change later
